# src/idcos.io/cloud-act2/scripts/profile.d/act2.sh
#!/bin/bash

export PATH=/usr/yunji/cloud-act2/bin:$PATH